# distribution/linux/scholarr.service
[Unit]
Description=Scholarr Academic Manager
After=network.target

[Service]
Type=simple
User=scholarr
Group=scholarr
EnvironmentFile=/etc/scholarr/scholarr.env
WorkingDirectory=/var/lib/scholarr
ExecStart=/opt/scholarr/venv/bin/uvicorn scholarr.app:create_app --factory --host 0.0.0.0 --port 8787
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
